# magicmemories.service
[Unit]
Description=Magic Memories Books - Flask Application
After=network.target mysql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/magic-memories
EnvironmentFile=/var/www/magic-memories/.env
Environment="PATH=/var/www/magic-memories/venv/bin:/usr/bin:/bin"
ExecStart=/var/www/magic-memories/venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 4 --timeout 300 --access-logfile /var/www/magic-memories/logs/access.log --error-logfile /var/www/magic-memories/logs/error.log app:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
